# nim-bot.service
[Unit]
Description=Nim Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/nim-bot
ExecStart=/usr/bin/node bot.js
Restart=always
RestartSec=5
EnvironmentFile=/home/ubuntu/nim-bot/.env

[Install]
WantedBy=default.target
